# s2t.service
[Unit]
Description=Speech to Text Conversion Service
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/defintek.io/s2t
ExecStart=/var/www/defintek.io/s2t/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=s2t
Environment="PATH=/var/www/defintek.io/s2t/.venv/bin"

[Install]
WantedBy=multi-user.target 